# systemd/pkgfile-update.service.in
[Unit]
Description=pkgfile database update
RequiresMountsFor=@DEFAULT_CACHEPATH@

[Service]
Type=oneshot
ExecStart=/usr/bin/pkgfile -u
Nice=19
StandardOutput=null
StandardError=journal
PrivateTmp=yes
PrivateDevices=yes
CapabilityBoundingSet=
NoNewPrivileges=yes
