# systemd/pkgfile-update.timer
[Unit]
Description=pkgfile database update timer

[Timer]
OnCalendar=daily
AccuracySec=6h
Persistent=yes

[Install]
WantedBy=mutli-user.target
